# .config/my-services/keyboard-disabler.service
[Unit]
Description=Keyboard Disabler Daemon
DefaultDependencies=no
After=local-fs.target
Before=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/keyboard-disabler
Restart=always
RestartSec=2
User=root

# Optional hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
